# www/config/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

/bin/busybox udhcpd
/etc/init.d/ntpdate start

exit 0
